# netemu.service
[Unit]
Description=NetEmu Rewrite
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/netemu/backend
ExecStart=/opt/netemu/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8080 --log-level info
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
